# build/ports/sssd/files/sssd.in
#!/bin/sh
#
# $FreeBSD: head/security/sssd/files/sssd.in 354743 2014-05-21 11:46:16Z feld $
#

# PROVIDE: sssd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
# XQUERY: -i "count(//kerberos/type[contains(., 'sss')]) and count(//kerberos/enable) > 0" -o "0" -b
# RCVAR: sssd

. /etc/rc.subr
. /etc/configxml.subr

name="sssd"
rcvar="sssd_enable"

# read configuration and set defaults
load_rc_config "$name"

: ${sssd_enable:=NO}
: ${sssd_conf="/usr/local/etc/sssd/sssd.conf"}
: ${sssd_flags="-f -D"}
sssd_config_file="/usr/local/etc/sssd/sssd.conf"

command="/usr/local/sbin/$name"
pidfile="/var/run/$name.pid"
required_files="${sssd_config_file}"
start_precmd="sssd_mkconf"

for i in db/sss db/sss_mc log/sssd run/sss/krb5.include.d run/sss/private run/sss; do
	if [ ! -d var/${i} ]; then mkdir -p /var/${i}; fi
done

# Check if required directories exists.
[ ! -d "/usr/local/etc/sssd/" ] && mkdir -p "/usr/local/etc/sssd/"

# Check if required config file exist.
[ ! -f "/usr/local/etc/sssd/sssd.conf" ] && touch "/usr/local/etc/sssd/sssd.conf"

chmod 600 /usr/local/etc/sssd/sssd.conf

sssd_mkconf()
{
	# Create sssd_config file
	/usr/local/bin/xml sel -t -m "//kerberos" \
		-m "sssdauxparam" -v "." -n -b \
		${configxml_file} | /usr/local/bin/xml unesc > ${sssd_config_file}
}

run_rc_command "$1"